# config/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=/usr/bin/uwsgi --ini /etc/uwsgi/uwsgi.ini --die-on-term
stdout_logfile=/var/log/uwsgi/app/access.log
stdout_logfile_maxbytes=1000000
stderr_logfile=/var/log/uwsgi/app/error.log
stderr_logfile_maxbytes=1000000

[program:nginx]
command=/usr/sbin/nginx
stdout_logfile=/var/log/nginx/access.log
stdout_logfile_maxbytes=1000000
stderr_logfile=/var/log/nginx/error.log
stderr_logfile_maxbytes=1000000
